# crontab
0 12 * * * echo "Writing daily report..." > /proc/1/fd/1 2>/proc/1/fd/2
0 12 * * * /code/main.py >> /var/log/cron.log 2>&1
